# root/etc/services.d/plexdrive/finish
#!/usr/bin/with-contenv bash

exec s6-setuidgid root \
  /bin/sh -c \
    "/bin/fusermount -u -z ${PLEXDRIVE_MOUNT_POINT}"
