# root/etc/services.d/rclone/run
#!/usr/bin/with-contenv bash

echo "# rclone waiting for unionfs_cleaner service"
s6-svwait ../unionfs_cleaner && sleep 1

echo "# rclone starting..."
exec s6-setuidgid plex /bin/bash -c '/start.sh'
